# onboarding/watchdog/cc-stall-watchdog.service
[Unit]
Description=Claude Code stall watchdog (dead-man's switch tick)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cc-stall-watchdog.sh
